# picast-control.timer
[Unit]
Description=Picast Controller Timer

[Timer]
OnCalendar=*:*:0/15

[Install]
WantedBy=timers.target
